# 3speak-ipfs-storage-management.service
[Unit]
Description=3Speak IPFS Storage Management Web Interface
Documentation=https://github.com/menobass/3speakstorageadmin
After=network.target mongodb.service

[Service]
Type=simple
# User and paths are set by install script - these are placeholders
User=__USER__
Group=__GROUP__
WorkingDirectory=__WORKING_DIR__
Environment=NODE_ENV=production
EnvironmentFile=__WORKING_DIR__/.env
ExecStart=/usr/bin/node __WORKING_DIR__/dist/web-server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=3speak-ipfs-storage-management

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=__WORKING_DIR__/logs
ProtectHome=read-only

[Install]
WantedBy=multi-user.target